# replay-sorcery.service.in
[Unit]
Description=ReplaySorcery: an open-source, instant-replay solution for Linux.

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/replay-sorcery
TimeoutStopSec=60
Restart=always

[Install]
WantedBy=graphical.target
